# jmxtrans/cron.daily/jmxtrans
#!/bin/sh

find /var/log/jmxtrans -type f -mtime +32 -name "jmxtrans-*.log.gz" -print0 | xargs -0 rm -f
